# crontab
*/5 * * * * python /home/soutys/TMProteins/branches/features/forms/TMProteins/manage.py runcrons > /home/soutys/cronlog
